# src/initrd/etc/systemd/system/generate_ssh_keys.service
[Unit]
Description=Generate SSH keys
Before=ssh.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate_ssh_host_keys
